# deploy/AsosFilter.conf
description "AsosFilter node process"

start on runlevel [2]
stop on runlevel [016]

console owner
chdir "/var/www/AsosFilter"
env NODE_ENV=production
exec node asosFilter.js >> /var/log/AsosFilter.log 2>&1
respawn
respawn limit 5 15
